# debian/openrobertalab.service
# /lib/systemd/system/openrobertalab.service

[Unit]
Description=OpenRoberta Lab connector
Conflicts=getty@tty2.service
After=systemd-user-sessions.service getty@tty2.service

[Service]
ExecStartPre=/bin/chown robot: /dev/tty2
ExecStart=/usr/bin/openrobertalab
ExecStopPost=/bin/chvt 1
Restart=always
User=robot
StandardInput=tty
TTYPath=/dev/tty2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=brickman.service
